# docker/auditwise.service
[Unit]
Description=AuditWise Production (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/auditwise
ExecStartPre=/bin/bash -c 'which docker'
ExecStartPre=/bin/bash -c 'docker compose -f docker-compose.vps.yml config >/dev/null'
ExecStart=/bin/bash -c 'docker compose -f docker-compose.vps.yml up -d --remove-orphans'
ExecStop=/bin/bash -c 'docker compose -f docker-compose.vps.yml stop'
ExecReload=/bin/bash -c 'docker compose -f docker-compose.vps.yml up -d --build --remove-orphans'
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
